# debian_pkg/myapp/debian/myapp.service
[Unit]
Description=myapp service

[Service]
ExecStart=/usr/bin/myapp
Restart=always

[Install]
WantedBy=multi-user.target