# debian_pkg/myapp/debian/postinst
#!/bin/sh

# start the service
systemctl start myapp
# Enable start on reboot
systemctl enable myapp