# debian_pkg/myapp/debian/prerm
#!/bin/sh

# stop the service
systemctl stop myapp

# Disable start on reboot
systemctl disable myapp